# third_party/bluetooth/bluetooth_default_connecter.service
[Unit]
Description=bluetooth default connecter
After=network.target bluetooth.target syslog.target

[Service]
User=ubuntu
ExecStart=/usr/bin/python /usr/sbin/bluetooth_default_connecter/start.py
ExecStop=echo "Done!"

[Install]
WantedBy=multi-user.target
